# ex04/99-usb-keyboard.rules
#   ACTION=="add"                   - デバイスが追加されたとき
#   SUBSYSTEM=="usb"                - USBサブシステム
#   ATTR{bInterfaceClass}=="03"     - HID (Human Interface Device) クラス
#   ATTR{bInterfaceSubClass}=="01"  - Boot Interface サブクラス
#   ATTR{bInterfaceProtocol}=="01"  - キーボードプロトコル

ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="03", ATTR{bInterfaceSubClass}=="01", ATTR{bInterfaceProtocol}=="01", RUN+="/sbin/modprobe main"

# usbhidからunbind (%k = デバイス名、例: 2-2:1.0)
ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="03", ATTR{bInterfaceSubClass}=="01", ATTR{bInterfaceProtocol}=="01", RUN+="/bin/sh -c 'echo %k > /sys/bus/usb/drivers/usbhid/unbind 2>/dev/null'"

# skel_driverにバインド
ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="03", ATTR{bInterfaceSubClass}=="01", ATTR{bInterfaceProtocol}=="01", RUN+="/bin/sh -c 'echo %k > /sys/bus/usb/drivers/skel_driver/bind 2>/dev/null'"

# ENV{INTERFACE}=="3/1/1" = HIDクラス(3)/Bootサブクラス(1)/キーボード(1)
ACTION=="remove", SUBSYSTEM=="usb", ENV{INTERFACE}=="3/1/1", RUN+="/sbin/modprobe -r main"
