# ex04/99-usb-keyboard.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/sbin/modprobe main"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/sbin/modprobe -r main"